# systemd/powertop.service
[Unit]
Description=Powertop auto-tune

[Service]
Type=idle
ExecStart=/usr/bin/powertop --auto-tune

[Install]
WantedBy=multi-user.target
