# _systemd/physlock.service
[Unit]
Description=Activate physlock on suspend

[Service]
User=evan
ExecStart=/usr/bin/physlock

[Install]
WantedBy=sleep.target
